# deploy/zenofewords.service
[Unit]
Description=zenofewords gunicorn
After=network.target postgresql.service

[Service]
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/zenofewords
EnvironmentFile=/home/deploy/zenofewords/.env
ExecStart=/home/deploy/zenofewords/.venv/bin/gunicorn zenofewords.wsgi:application --bind 127.0.0.1:8001 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
